# debian/distiller-update.timer
[Unit]
Description=Distiller Update Notifier Timer
Documentation=https://github.com/pamir-ai-pkgs/distiller-update

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=10min
# Make the timer persistent
Persistent=true

[Install]
WantedBy=timers.target
