# debian/postinst
#!/bin/sh
set -e

export PATH="/root/.local/bin:/root/.cargo/bin:$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:/usr/bin:/bin:$PATH"

case "$1" in
configure)
	cd /opt/distiller-update

	echo "INFO: Setting up virtual environment with uv..."

	# Setup virtual environment
	[ -d ".venv" ] && rm -rf ".venv"
	if ! uv venv --system-site-packages 2>/dev/null; then
		echo "INFO: Creating venv without system-site-packages..."
		uv venv || {
			echo "WARNING: Failed to create virtual environment with uv"
			exit 1
		}
	fi

	# Install dependencies
	uv sync

	VENV_PATH="/opt/distiller-update/.venv"

	# Set permissions
	find /opt/distiller-update -type d -exec chmod 755 {} \; \
		-o -type f -name "*.py" -exec chmod 644 {} \;
	chmod +x /usr/bin/distiller-update

	# Update systemd service files to use venv Python
	if [ -f /etc/systemd/system/distiller-update.service ]; then
		sed -i "s|/usr/bin/python3|$VENV_PATH/bin/python|g" /etc/systemd/system/distiller-update.service
		sed -i "s|ExecStart=distiller-update|ExecStart=/usr/bin/distiller-update|g" /etc/systemd/system/distiller-update.service
	fi

	# Ensure cache directory exists with proper permissions
	mkdir -p /var/cache/distiller-update
	chmod 755 /var/cache/distiller-update

	# Reload systemd if not in build environment
	if [ -f /boot/armbianEnv.txt ] || [ -f /boot/firmware/config.txt ]; then
		systemctl daemon-reload 2>/dev/null || true
		systemctl enable distiller-update.timer 2>/dev/null || true
		echo "INFO: Distiller Update service installed."
	fi

	# Migration to testing distribution (v2.3.0+)
	echo ""
	echo "=== Pamir AI Repository Migration ==="
	echo "This package will automatically migrate your system from"
	echo "unstable to testing distribution and install new packages."
	echo ""

	# Install migration scripts
	install -d /usr/local/bin
	install -m 0755 /opt/distiller-update/scripts/migrate-to-testing.sh \
		/usr/local/bin/migrate-to-testing.sh
	install -m 0755 /opt/distiller-update/scripts/platform-detect.sh \
		/usr/local/bin/platform-detect.sh

	# Create systemd service
	cat >/etc/systemd/system/pamir-migration.service <<'MIGRATION_EOF'
[Unit]
Description=Pamir AI Repository Migration Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/migrate-to-testing.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
MIGRATION_EOF

	# Enable and start service (non-blocking)
	systemctl daemon-reload
	systemctl enable pamir-migration.service 2>/dev/null || true
	systemctl start pamir-migration.service --no-block || true

	echo "Migration service started in background"
	echo "Monitor progress: journalctl -u pamir-migration.service -f"
	echo "Check log: tail -f /var/log/pamir-migration.log"
	echo ""
	echo "Migration will complete automatically in 5-15 minutes"
	echo "======================================================="
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "ERROR: postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
