# debian/postinst
#!/bin/sh
set -e

export PATH="/root/.local/bin:/root/.cargo/bin:$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:$PATH"

case "$1" in
configure)
	cd /opt/distiller-update

	# Setup virtual environment
	[ -d ".venv" ] && rm -rf ".venv"
	if ! uv venv --system-site-packages 2>/dev/null; then
		echo "INFO: Creating venv without system-site-packages..."
		uv venv || {
			echo "ERROR: Failed to create virtual environment with uv"
			exit 1
		}
	fi

	# Install dependencies
	echo "INFO: Installing dependencies with uv sync..."
	uv sync

	# Set permissions
	find /opt/distiller-update -type d -exec chmod 755 {} \; \
		-o -type f -name "*.py" -exec chmod 644 {} \;
	chmod +x /usr/bin/distiller-update

	# Ensure cache directory exists with proper permissions
	mkdir -p /var/cache/distiller-update
	chmod 755 /var/cache/distiller-update

	# Reload systemd
	systemctl daemon-reload 2>/dev/null || true
	systemctl enable distiller-update.timer 2>/dev/null || true

	echo "INFO: distiller-update installed successfully."
	echo "INFO: Start service: sudo systemctl start distiller-update.timer"
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
