# debian/postinst
#!/bin/sh
set -e

export PATH="/root/.local/bin:/root/.cargo/bin:$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:/usr/bin:/bin:$PATH"

case "$1" in
configure)
	cd /opt/distiller-update

	echo "INFO: Setting up virtual environment with uv..."

	# Setup virtual environment
	[ -d ".venv" ] && rm -rf ".venv"
	if ! uv venv --system-site-packages 2>/dev/null; then
		echo "INFO: Creating venv without system-site-packages..."
		uv venv || {
			echo "WARNING: Failed to create virtual environment with uv"
			exit 1
		}
	fi

	# Install dependencies
	uv sync

	VENV_PATH="/opt/distiller-update/.venv"

	# Set permissions
	find /opt/distiller-update -type d -exec chmod 755 {} \; \
		-o -type f -name "*.py" -exec chmod 644 {} \;
	chmod +x /usr/bin/distiller-update

	# Update systemd service files to use venv Python
	if [ -f /etc/systemd/system/distiller-update.service ]; then
		sed -i "s|/usr/bin/python3|$VENV_PATH/bin/python|g" /etc/systemd/system/distiller-update.service
		sed -i "s|ExecStart=distiller-update|ExecStart=/usr/bin/distiller-update|g" /etc/systemd/system/distiller-update.service
	fi

	# Ensure cache directory exists with proper permissions
	mkdir -p /var/cache/distiller-update
	chmod 755 /var/cache/distiller-update

	# Reload systemd if not in build environment
	if [ -f /boot/armbianEnv.txt ] || [ -f /boot/firmware/config.txt ]; then
		systemctl daemon-reload 2>/dev/null || true
		systemctl enable distiller-update.timer 2>/dev/null || true
		echo "INFO: Distiller Update service installed."
	fi
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "ERROR: postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
