# debian/prerm
#!/bin/bash
# Pre-removal script for distiller-update

set -e

case "$1" in
remove | upgrade | deconfigure)
	# Stop and disable the timer
	if [ -d /run/systemd/system ]; then
		# Stop the timer
		systemctl stop distiller-update.timer || true

		# Stop the service if it's running
		systemctl stop distiller-update.service || true

		# Disable the timer
		systemctl disable distiller-update.timer || true

		# Stop and disable migration service if it exists
		systemctl stop pamir-migration.service 2>/dev/null || true
		systemctl disable pamir-migration.service 2>/dev/null || true
	fi

	# Clean up migration artifacts
	rm -f /etc/systemd/system/pamir-migration.service
	rm -f /usr/local/bin/migrate-to-testing.sh
	rm -f /usr/local/bin/platform-detect.sh
	systemctl daemon-reload 2>/dev/null || true

	# Clean up MOTD file if it exists
	if [ -f /etc/update-motd.d/99-distiller-updates ]; then
		rm -f /etc/update-motd.d/99-distiller-updates
	fi
	;;

failed-upgrade) ;;

*)
	echo "ERROR: prerm called with unknown argument '$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
